# claude-glm-proxy.service
[Unit]
Description=Claude GLM Proxy - Proxy Anthropic Haiku requests to Z.ai GLM
After=network.target

[Service]
Type=simple
User=<USER>
WorkingDirectory=/home/<USER>/claude-glm-proxy
ExecStart=/home/<USER>/.bun/bin/bun run proxy.mjs
Restart=always
RestartSec=3
Environment="PATH=/home/<USER>/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"

# Logging
StandardOutput=append:/home/<USER>/claude-glm-proxy/logs/proxy.log
StandardError=append:/home/<USER>/claude-glm-proxy/logs/proxy-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
